# scream.service
[Unit]
Description=Scream
After=graphical.target

[Service]
ExecStartPre=/bin/sleep 15
Type=simple
ExecStart=/bin/scream 